# service/vanguards.service
[Unit]
Description=Vanguards Tor guard protection daemon
After=network.target tor.service
Wants=tor.service

[Service]
Type=simple
ExecStart=vanguards --logfile /var/log/vanguards.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
