# service/termux/run
#!/data/data/com.termux/files/usr/bin/sh
exec vanguards --logfile /data/data/com.termux/files/usr/var/log/vanguards.log
